# etc/crontab
0 8-18 * * 1-5 /usr/local/bin/reminder-upto
